# duomus/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community App: DuoMus
# Configures the nginx service
# ==============================================================================
readonly NGINX_CONFIG_PATH=/etc/nginx/http.d/default.conf
declare ipaddress

bashio::log.info "Starting nginx configuration update..."

# Check if a ipaddress is set
ipaddress=$(bashio::config 'ipaddress')
ipaddress=$(bashio::string.lower "${ipaddress}")

if ! bashio::var.has_value "${ipaddress}"; then
    bashio::log.info 'No ip address is defined in the configuration!'
    exit 0
fi

# Update the proxy url.
sed -i "s/proxy_pass .*/proxy_pass http:\/\/${ipaddress}\/;/" "${NGINX_CONFIG_PATH}" \
    || bashio::exit.nok 'Failed configuring ip address'

bashio::log.notice 'Changed the proxy to: http://'${ipaddress}
